# systemd/user/systemreboot.service
[Unit]
Description=Jack Audio Server

[Service]
Type=simple
# ExecStartPre=/bin/sleep 30
ExecStart=/sbin/shutdown -r now