# EN325/wonder/twonder/riviera/systemd/twonder.target
[Unit]
Description=Twonders
Requires=twonder@0.service twonder@1.service twonder@2.service twonder@3.service

[Install]
WantedBy=default.target