# EN325/ambisonics/systemd/allra_decoder.service
[Unit]
Description=AllRADecoder
After=jack.service

[Install]
WantedBy=default.target

[Service]
Type=simple
ExecStart=AllRADecoder
LimitRTPRIO=95
LimitRTTIME=infinity
LimitMEMLOCK=infinity